# .config/systemd/user/nm-applet.service
[Unit]
Description=nm-applet docklet
After=systemd-user-sessions.service
After=network.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=default.target
