# .config/systemd/user/swayidle.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/home/tcrha/.config/systemd/user/swayidle.sh
Restart=on-failure

[Install]
WantedBy=default.target
