# .config/systemd/user/wl-paste.service
[Unit]
Description=Watch clipboard and store in clipman
After=graphical-session.target

[Service]
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
ExecStart=/usr/bin/wl-paste -t text --watch clipman store

[Install]
WantedBy=default.target
