# .config/systemd/user/wlsunset.service
[Unit]
Description=wlsunset dim screen at sunset
After=sway-session.target
After=systemd-user-sessions.service

[Service]
Type=simple
StandardOutput=append:/tmp/wlsunset.log
StandardError=append:/tmp/wlsunset.log
Restart=always
ExecStart=/home/tcrha/.config/systemd/user/wlsunset-start.sh

[Install]
WantedBy=default.target
