# .config/systemd/user/wlsunset.service
[Unit]
Description=wlsunset dim screen at sunset
After=systemd-user-sessions.service
After=network.target
After=network-online.target

[Service]
Type=simple
StandardOutput=append:/tmp/wlsunset.log
StandardError=append:/tmp/wlsunset.log
Restart=always
ExecStart=/bin/bash -c 'rm /tmp/wlsunset.state || true && echo "{\\"text\\":\\"\\", \\"tooltip\\":\\"hello\\"}" > /tmp/wlsunset.state && /usr/bin/wlsunset -l -33.7 -L 150.3'

[Install]
WantedBy=default.target
