# udev/etc_udev_rules.d/61-onacpower.rules
# Rule for when switching to AC Power
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/local/bin/set-power-profile performance"
